# groupalarm-auto-update.service
[Unit]
Description=GroupAlarm Monitor Auto-Update Job v5.6
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Runs as the kiosk user so it can pause/resume Chromium in the active X session.
# Privileged steps (apt-get, systemctl, shutdown) go through the sudoers allow-list.
User=groupalarm
Group=groupalarm
Environment="HOME=/home/groupalarm"
WorkingDirectory=/usr/local/bin
ExecStart=/usr/local/bin/groupalarm-auto-update.sh
StandardOutput=journal
StandardError=journal
